# src/lib/pihole-update-firmware.service
[Unit]
Description=Manage PI-hole by updating firmware

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pihole-manager update --firmware