# src/lib/pihole-log-flush.service
[Unit]
Description=Flush Pihole log

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pihole-manager flush