# src/lib/pihole-update-gravity.service
[Unit]
Description=Manage PI-hole by updating the gravity

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pihole-manager update --gravity